# guardian-monitor.timer
[Unit]
Description=Guardian Monitor Timer - 每5分钟执行一次健康检查
Requires=guardian-monitor.service

[Timer]
OnBootSec=1min
OnUnitActiveSec=5min
Unit=guardian-monitor.service

[Install]
WantedBy=timers.target
