# oopsclaw-monitor.service
[Unit]
Description=Oopsclaw Monitor Service
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/home/duanz/.oopsclaw/monitor.sh
WorkingDirectory=/home/duanz/.oopsclaw
User=duanz

[Install]
WantedBy=default.target
