# oopsclaw-event.service
[Unit]
Description=Oopsclaw Event - OpenClaw故障自动修复
After=network-online.target

[Service]
Type=oneshot
ExecStart=/bin/bash /home/duanz/.oopsclaw/oopsclaw-event.sh
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
